# KoboRoot/etc/udev/rules.d/update_conf.rules
KERNEL=="loop0", RUN+="/root/update_conf.sh"
